# apcups/apcupsd_exporter.service
[Unit]
Description=APC UPS Power Control Monitor for Prometheus

[Service]
ExecStart=/usr/bin/apcupsd_exporter

[Install]
WantedBy=multi-user.target
